# deploy/design-audit.service
[Unit]
Description=design-os weekly watchlist sweep

[Service]
Type=oneshot
EnvironmentFile=%h/.config/design-os/env
ExecStart=%h/design-os/.venv/bin/design-os --watchlist %h/design-os/watchlist.yaml
WorkingDirectory=%h/design-os
